# example/journalcheck.service
[Unit]
Description=journalcheck

[Service]
Type=oneshot
User=lynix
Environment=JC_LOGLEVEL=0..4
ExecStart=/usr/bin/checkrun journalcheck
