# example/journalcheck.timer
[Unit]
Description=Trigger journalcheck

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
